# deploy/systemd/project-steele-update.timer
[Unit]
Description=Check GitHub main for Project Steele updates

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
